# test-e2e/resources/dummy.service
[Unit]
Description=Touches a file, demonstrating that the service ran to completion

[Service]
ExecStart=touch /tmp/my-timer-was-here
